# proxy/com.mewvault.chromadb.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mewvault.chromadb</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/Mohabbat/Jan/mewvault/venv/bin/chroma</string>
        <string>run</string>
        <string>--path</string>
        <string>/Users/Mohabbat/.mewvault/chromadb</string>
        <string>--port</string>
        <string>8001</string>
        <string>--host</string>
        <string>localhost</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/Mohabbat/Jan/mewvault/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/Mohabbat</string>
        <key>ANONYMIZED_TELEMETRY</key>
        <string>False</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/Mohabbat/.mewvault/logs/chromadb.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/Mohabbat/.mewvault/logs/chromadb.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/Mohabbat/Jan/mewvault</string>
</dict>
</plist>
